# misc/hdapsd@.service.in
[Unit]
Description=hdapsd hard drive active protection system daemon - %I
Documentation=man:hdapsd(8)

[Service]
SyslogIdentifier=%p(%I)
Nice=-5
ExecStart=@sbindir@/hdapsd -a -d %I
Restart=on-abort
