# misc/hdapsd.service.in
[Unit]
Description=hdapsd hard drive active protection system daemon
Documentation=man:hdapsd(8)

[Service]
SyslogIdentifier=%p
Nice=-5
ExecStart=@sbindir@/hdapsd --syslog
Restart=on-abort

[Install]
WantedBy=multi-user.target
